# scripts/crontab/dcron
@reboot sx /opt/distri/scripts/crontab/runner/back_end.sh &
@reboot sxh /opt/distri/scripts/crontab/runner/front_end.sh &
@reboot sxh /opt/distri/scripts/crontab/runner/printer.sh &
@reboot date >> ~/Escritorio/jojo.txt &
