# scripts/crontab/dcron
@reboot sleep 200 && sh /opt/distri/scripts/crontab/_init.sh
@reboot sleep 210 && sh /opt/distri/scripts/crontab/updater.sh

@reboot sleep 220 && sh /opt/distri/scripts/crontab/runner/back_end.sh
@reboot sleep 230 && sh /opt/distri/scripts/crontab/runner/front_end.sh
@reboot sleep 240 && sh /opt/distri/scripts/crontab/runner/printer.sh

@reboot sleep 45 && date >> ~/Escritorio/jjx.txt
@reboot sleep 50 && date >> ~/Escritorio/yjj.txt
