# scripts/distri.service
[Unit]
Description=Distri Boot

[Service]
ExecStart=/opt/distri/scripts/booting.sh start

WantedBy=multi-user.target
[Install]

